# frontend/appflowy_web_app/supervisord.conf
[supervisord]
nodaemon=true

[program:start]
command=sh /app/start.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/start.err.log
stdout_logfile=/var/log/start.out.log
