# frontend/appflowy_web_app/deploy/supervisord.conf
[supervisord]
nodaemon=true

[program:bun]
command=sh /app/start.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/bun.err.log
stdout_logfile=/var/log/bun.out.log
